# supervisor/dropbear.conf
[program:dropbear]
command=/usr/sbin/dropbear -F -p 22
autostart=true
autorestart=true
redirect_stderr=true
