# supervisor/node-gpio.conf
[program:node-gpio]
command=node /app/gpio.js
autostart=true
autorestart=true
redirect_stderr=true
