# supervisor/node.conf
[program:npm]
command=node /app/adc.js
autostart=true
autorestart=true
redirect_stderr=true
