# supervisor/node.conf
[program:node-main]
command=node /app/adc.js
autostart=true
autorestart=true
redirect_stderr=true
startsecs=5
