# supervisor/owfs.conf
[program:owfs]
command=/usr/bin/owfs --foreground
autostart=true
autorestart=true
redirect_stderr=true
