# supervisor/telldusd.conf
[program:telldusd]
command=/usr/sbin/telldusd --nodaemon
autostart=true
autorestart=true
redirect_stderr=true
